# home/.config/systemd/user/fusuma.service
[Unit]
Description=Fusuma touchpad gestures

[Service]
ExecStart=/bin/fusuma -v
Restart=on-failure
RestartSec=10s

